# roles/music-setup/files/aspp-music-finalcountdown.timer
[Unit]
Description=Pushit

[Timer]
OnCalendar=Saturday 17:20
AccuracySec=1us

[Install]
WantedBy=timers.target

