# mkosi/mkosi.extra/usr/lib/systemd/system/growfs-root.service
[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/systemd-growfs /

[Install]
WantedBy=sysinit.target
