# mkosi/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/grow-root.conf
[Service]
ExecStartPost=/usr/lib/systemd/systemd-growfs /
